# .config/autostart/set-kb-layout.sh.desktop
[Desktop Entry]
Exec=/home/matus/set-kb-layout.sh
Icon=dialog-scripts
Name=set-kb-layout.sh
Path=
Type=Application
X-KDE-AutostartScript=true
